# zdclient-locale-1.2.2/luci-zdclient/root/etc/uci-defaults/luci-zdclient
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@zdclient[-1]
	add ucitrack zdclient
	set ucitrack.@zdclient[-1].init=zdclient
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
